# tools/install/veska-daemon.service.tmpl
[Unit]
Description=Veska code intelligence daemon
After=network.target

[Service]
ExecStart={{.BinaryPath}}
Restart=on-failure
RestartSec=5
Environment=VESKA_HOME={{.VeskaHome}}
StandardOutput=append:{{.VeskaHome}}/daemon.log
StandardError=append:{{.VeskaHome}}/daemon.log

[Install]
WantedBy=default.target
